# packages/os/overlay/etc/systemd/system/umbrel.service
[Unit]
Description=Umbrel daemon
Wants=umbrel-custom-pre-start.service
After=network-online.target docker.service umbrel-custom-pre-start.service

[Service]
TimeoutStopSec=15min
ExecStart=umbreld --data-directory=/home/umbrel/umbrel
Restart=always
# This prevents us hitting restart rate limits and ensures we keep restarting
# indefinitely.
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
