# packages/os/overlay-common/etc/systemd/system/umbrel-tty-message.service
[Unit]
Description=Display Umbrel access information on TTY
After=umbrel.service

[Service]
ExecStart=/opt/umbrel-tty-message/umbrel-tty-message
Type=oneshot

[Install]
WantedBy=multi-user.target