# packages/os/overlay-amd64/etc/systemd/system/umbrel-expand-partitions.service
[Unit]
Description=Repartition Root Disk
DefaultDependencies=no
Conflicts=shutdown.target
After=initrd-usr-fs.target
Before=initrd-root-fs.target shutdown.target data.mount systemd-growfs@data.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/systemd-repart --dry-run=no

[Install]
WantedBy=data.mount